# mediaserver.service
[Unit]
Description=Media Sync Server
After=network.target

[Service]
Type=simple
# Замени YOUR_USER на своё имя пользователя
User=YOUR_USER
ExecStart=/usr/local/bin/mediaserver /etc/mediaserver/config.toml
Restart=on-failure
RestartSec=5
Environment=RUST_LOG=info

# Ограничения безопасности (опционально, но правильно)
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/mnt/media /tmp/mediaserver

[Install]
WantedBy=multi-user.target
